# infra/arc402-workroom.service
[Unit]
Description=ARC-402 Workroom
After=network-online.target docker.service
Wants=network-online.target
# Restart workroom if it crashes; also restart on success (container died)
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
# NVM-managed node — must include full PATH
Environment=PATH=/home/lego/.nvm/versions/node/v22.22.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/home/lego

# arc402 workroom start exits after launching the container
# Use a wrapper that blocks: start, then tail container logs (keeps unit alive while container runs)
ExecStart=/bin/bash -c '\
  /home/lego/.nvm/versions/node/v22.22.0/bin/arc402 workroom start && \
  docker wait arc402-workroom; \
  exit 1'

# On stop: stop the workroom gracefully
ExecStop=/home/lego/.nvm/versions/node/v22.22.0/bin/arc402 workroom stop

Restart=always
RestartSec=10s
TimeoutStartSec=120
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=arc402-workroom

[Install]
WantedBy=default.target
